# files/etc/network/if-up.d/rlbuild-bidon
#!/bin/sh

# This makes it easier to re-run 'puppet apply'
ip -6 addr show dev eth1 | grep -q 2607:f2c0:f00f:2911::10 || {
  ip -6 addr add 2607:f2c0:f00f:2911::10 dev eth1
}
